# djangopostgresdemo/supervisor.conf
[supervisord]
logfile=%(ENV_HOME)s/supervisor/supervisord.log
pidfile=%(ENV_HOME)s/supervisor/supervisord.pid

[unix_http_server]
file=%(ENV_HOME)s/supervisor/supervisor.sock

[supervisorctl]
serverurl=unix://%(ENV_HOME)s/supervisor/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# Listener - AUTO-RESTART on crash
[program:django-listener]
command=%(ENV_DJANGO_PROJECT)s/venv/bin/python manage.py listen_submissions
directory=%(ENV_DJANGO_PROJECT)s
autostart=true
autorestart=true
startsecs=5
stopasgroup=true
stdout_logfile=%(ENV_HOME)s/supervisor/listener.log
stderr_logfile=%(ENV_HOME)s/supervisor/listener_error.log
loglevel=info
environment=DJANGO_SETTINGS_MODULE=dbapp.settings

# Celery Worker 1 - AUTO-RESTART
[program:celery-worker-1]
command=%(ENV_DJANGO_PROJECT)s/venv/bin/celery -A dbapp worker --loglevel=info --hostname=worker1@%%h
directory=%(ENV_DJANGO_PROJECT)s
autostart=true
autorestart=true
startsecs=10
stopasgroup=true
stdout_logfile=%(ENV_HOME)s/supervisor/celery_worker1.log
stderr_logfile=%(ENV_HOME)s/supervisor/celery_worker1_error.log
loglevel=info
environment=DJANGO_SETTINGS_MODULE=dbapp.settings

# Celery Worker 2 - AUTO-RESTART
[program:celery-worker-2]
command=%(ENV_DJANGO_PROJECT)s/venv/bin/celery -A dbapp worker --loglevel=info --hostname=worker2@%%h
directory=%(ENV_DJANGO_PROJECT)s
autostart=true
autorestart=true
startsecs=10
stopasgroup=true
stdout_logfile=%(ENV_HOME)s/supervisor/celery_worker2.log
stderr_logfile=%(ENV_HOME)s/supervisor/celery_worker2_error.log
loglevel=info
environment=DJANGO_SETTINGS_MODULE=dbapp.settings

# Celery Beat (scheduler)
[program:celery-beat]
command=%(ENV_DJANGO_PROJECT)s/venv/bin/celery -A dbapp beat --loglevel=info
directory=%(ENV_DJANGO_PROJECT)s
autostart=true
autorestart=true
startsecs=10
stdout_logfile=%(ENV_HOME)s/supervisor/celery_beat.log
stderr_logfile=%(ENV_HOME)s/supervisor/celery_beat_error.log
loglevel=info
environment=DJANGO_SETTINGS_MODULE=dbapp.settings

# Celery Flower (monitoring UI)
[program:celery-flower]
command=%(ENV_DJANGO_PROJECT)s/venv/bin/celery -A dbapp flower --port=5555 --loglevel=info
directory=%(ENV_DJANGO_PROJECT)s
autostart=true
autorestart=true
startsecs=10
stdout_logfile=%(ENV_HOME)s/supervisor/celery_flower.log
stderr_logfile=%(ENV_HOME)s/supervisor/celery_flower_error.log
loglevel=info
environment=DJANGO_SETTINGS_MODULE=dbapp.settings

[group:celery]
programs=celery-worker-1,celery-worker-2,celery-beat,celery-flower
